# dists/systemd/scalix-monitoring-web.service
[Unit]
Description=Scalix Monitoring Web
